# debian/postinst
#!/bin/sh
update-rc.d shm2mqtt defaults
invoke-rc.d shm2mqtt start
